# crontab
* * * * * /project/scripts/run_pipeline.sh >> /var/log/cron.log 2>&1
